# debian/ccnq3-traces.cron.hourly
#!/bin/sh
SRC="`ccnq3 get_config_source`"
exec "${SRC}/applications/traces/cron.hourly"
