# debian/ccnq3.cron.daily
#!/bin/sh

NAME=ccnq3
USER=ccnq3
WORKDIR=/opt/${NAME}

(cd ${WORKDIR}/src && su -s /bin/sh -c "npm run-script compact" $USER)
