# debian/ccnq3-traces.cron.hourly
#!/bin/sh
exec /opt/ccnq3/src/applications/traces/cron.hourly
