# debian/ccnq3.cron.daily
#!/bin/sh

NAME=ccnq3
USER=ccnq3
WORKDIR=/opt/${NAME}

su -s /bin/sh -c "(cd ${WORKDIR}/src && npm run-script compact)" $USER
